# helper/restic@.service
[Unit]
Description=Restic Backup Service for %i

[Service]
Type=simple
User=root
ExecStart=/bin/restic backup %i
