# helper/systemd.service
[Unit]
Description=Restic Backup Service for %i

[Service]
Type=simple
User=root
ExecStart=/bin/restic-cli backup %i
